# .config/systemd/user/AdGuardHome.service
[Unit]
Description=AdGuardHome

[Service]
ExecStart=%h/apps/AdGuardHome/AdGuardHome

[Install]
WantedBy=default.target

