# .config/systemd/user/backupDotfiles.timer
[Unit]
Description=Daily backup of dotfiles to GitHub

[Timer]
OnCalendar=*-*-* 18:30:00
Persistent=true

[Install]
WantedBy=timers.target
