# .config/systemd/user/clearBackups.service
[Unit]
Description=Delete vim backup files (*~)

[Service]
Type=simple
ExecStart=/usr/bin/find /home/kchuangk -name *~ -type f -delete
