# .config/systemd/user/clrSwaps.service
[Unit]
Description=Delete vim backup files (.*.swp)

[Service]
Type=simple
ExecStart=/usr/bin/find %h -name '.*.swp' -type f -delete
