# .config/systemd/user/gitbackup.service
[Unit]
Description=Timer to run GIT for Documents directory daily

[Service]
Type=simple
ExecStart=$HOME/Documents/.git/gitbackup.sh
