# .config/systemd/user/gitbackup.timer
[Unit]
Description=Run daily git backup automatically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
