# .config/systemd/user/latex_all.service
[Unit]
Description=Compile all latex resume files

[Service]
Type=simple
ExecStart=%h/Documents/resumes/latex_all.sh

