# .config/systemd/user/pCloudBackup.timer
[Unit]
Description=pCloud Backup

[Timer]
OnCalendar=*-*-* 12:30:00
Persistent=true

[Install]
WantedBy=timers.target
