# .config/systemd/user/pCloudBackup.service
[Unit]
Description=pCloud Backup

[Service]
Type=simple
ExecStart=/opt/borg_backup.sh

