# .config/systemd/user/pCloudBackup.service
[Unit]
Description=pCloud Backup

[Service]
Type=simple
ExecStart=%h/scripts/pCloud_borg_backup.sh

