# .config/systemd/user/pCloudBackup.timer
[Unit]
Description=pCloud Backup

[Timer]
OnBootSec = 0min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
